# dotfiles/physical/xfce/autostart/xrandr.desktop
[Desktop Entry]
Name=Xrandr
Comment=Xrandr settings
Exec=sh -c "xrandr --output Virtual-1 --mode 2048x1152"
Terminal=false
Type=Application